# deploy/supervisor/homework-api.conf
; Homework API（ABP HttpApi.Host）—— supervisor 托管
; 放到 /etc/supervisor/conf.d/ 下：
;   sudo cp deploy/supervisor/homework-api.conf /etc/supervisor/conf.d/
;   sudo supervisorctl reread && sudo supervisorctl update && sudo supervisorctl start homework-api

[program:homework-api]
directory=/var/www/homework/api
command=/usr/bin/dotnet /var/www/homework/api/Homework.HttpApi.Host.dll
autostart=true
autorestart=true
startretries=3
startsecs=10
stopsignal=TERM
stopwaitsecs=25
user=www-data
; 关键环境变量：
;  - ASPNETCORE_URLS 仅监听本机回环，由 nginx 反代
;  - DOTNET_gcServer=0 用 Workstation GC，省内存（2H4G 友好）
;  - 环境为 Production：不加载 appsettings.Development.json（PlayDemo 种子/dev 日志不会触发）
environment=ASPNETCORE_ENVIRONMENT="Production",ASPNETCORE_URLS="http://127.0.0.1:5000",DOTNET_gcServer="0",DOTNET_CLI_TELEMETRY_OPTOUT="1"
stdout_logfile=/var/log/homework/api.out.log
stderr_logfile=/var/log/homework/api.err.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
